# docker/build/supervisord-server-symfony.conf
[program:symfonyserver]
command=/run_server.sh
numprocs=1
autostart=true
autorestart=true